# deploy/macos/com.jedisherpa.wizardjoeavatar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jedisherpa.wizardjoeavatar</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/.venv/bin/python</string>
    <string>__ROOT__/tools/run_wizard_avatar_server.py</string>
    <string>--port</string>
    <string>8765</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>3</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/WizardJoeAvatar.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/WizardJoeAvatar.error.log</string>
</dict>
</plist>
